# supervisord.conf
[supervisord]
nodaemon=true

[program:start]
command=/opt/startup/start.sh
stdout_logfile=/var/log/supervisor/start.log
stderr_logfile=/var/log/supervisor/start.err
autorestart=true
